# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:go]
command=./server

[program:nginx]
command=nginx -g "daemon off;"